# leash.service
[Unit]
Description=Leash NDI Source Control
After=network.target
StartLimitAction=none

[Service]
Type=simple
User=YOUR_USER
Group=YOUR_USER
WorkingDirectory=/home/YOUR_USER/Leash
EnvironmentFile=/home/YOUR_USER/Leash/.env
ExecStart=/home/YOUR_USER/Leash/venv/bin/gunicorn \
    --workers 1 \
    --threads 4 \
    --bind 0.0.0.0:5000 \
    --timeout 120 \
    --access-logfile - \
    --error-logfile - \
    run:app
Restart=always
RestartSec=10s
StandardOutput=journal
StandardError=journal
SyslogIdentifier=leash
StartLimitIntervalSec=300
StartLimitBurst=10

[Install]
WantedBy=multi-user.target
